# .ebextensions/crontab
0 * * * * /usr/local/bin/node /var/app/current/index.js
